# sway/.config/systemd/user/wmenu-update-cache.service
[Unit]
Description=Update cache for wmenu-desktop-entries

[Service]
Type=oneshot
ExecStart=%h/.local/bin/wmenu-create-cache.sh
