# sway-common/.config/systemd/user/wmenu-update-cache.timer
[Unit]
Description=Update cache for wmenu-desktop-entries

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
