# core/.config/systemd/user/onedrive-sync.service
[Unit]
Description=Synchronize OneDrive directory

[Service]
Type=oneshot
ExecStart=/usr/bin/onedrive --sync
