# core/.config/systemd/user/tldr-update.service
[Unit]
Description=Update tldr's cache

[Service]
Type=oneshot
ExecStart=/usr/bin/tldr --update
