# core/.config/systemd/user/tldr-update.timer
[Unit]
Description=Update tldr's cache

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
